# etc/supervisord.vusion.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9010        ; (ip_address:port specifier, *:port for all iface)


[supervisord]
logfile=./tmp/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=debug               ; (log level;default info; others: debug,warn,trace)
pidfile=./tmp/pids/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://localhost:9010 ; use an http:// url to specify an inet socket

[program:multi_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=vusion.VusionMultiWorker
    --config=./etc/ttc_multiworker.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:export_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=vusion.ExportWorker
    --config=./etc/ttc_export_worker.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:stats_worker]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=vusion.StatsWorker
    --config=./etc/ttc_stats_worker.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:dispatcher]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=dispatchers.DynamicDispatchWorker
    --config=./etc/ttc_dispatcher.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:dispatcher_256_8181]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=vumi.dispatchers.BaseDispatchWorker
    --config=./etc/ttc_dispatcher_256_8181.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:dispatcher_256_8181_yo]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=vumi.dispatchers.BaseDispatchWorker
    --config=./etc/ttc_dispatcher_256_8181_yo.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:ttc_garbage]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=vusion.GarbageWorker
    --config=./etc/ttc_garbage_worker.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:transport_256_8282_yo]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=transports.YoHttpTransport
    --config=./etc/ttc_256_8282_yo_http.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:transport_256_8181_yo]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=transports.YoHttpTransport
    --config=./etc/ttc_256_8181_yo_http.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:transport_256_8181_mtn]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=transports.YoHttpTransport
    --config=./etc/ttc_256_8181_mtn_http.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:transport_256_8181_yo_bulk]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    start_worker
    --vumi-config=./etc/vumi_config.yaml
    --worker-class=transports.YoHttpTransport
    --config=./etc/ttc_256_8181_yo_http_bulk.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:transport_forward_http]
numprocs=1
numprocs_start=1
process_name=%(program_name)s_%(process_num)s
command=twistd -n
    --pidfile=./tmp/pids/%(program_name)s_%(process_num)s.pid
    --log=DEBUG
    start_worker
    --worker-class=transports.ForwardHttp
    --vumi-config=./etc/vumi_config.yaml
    --config=./etc/ttc_forward_http.yaml
stdout_logfile=./tmp/log/%(program_name)s_%(process_num)s.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=./tmp/log/%(program_name)s_%(process_num)s.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
